# debian/salor-retail.cron.d
#
# Regular cron jobs for the salor-hospitality package
#
0 2 * * * /usr/bin/salor-maintainance r
